# resources/rootfs/overlay/etc/systemd/system/fcnet.service
[Service]
Type=oneshot
ExecStartPre=/usr/bin/udevadm settle
ExecStart=/usr/local/bin/fcnet-setup.sh
[Install]
WantedBy=sshd.service
